# gateway/com.kitty.model-digest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kitty.model-digest</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/jacobbrizinski/Projects/kitty/venv/bin/python</string>
        <string>-c</string>
        <string>import sys; sys.path.insert(0, '/Users/jacobbrizinski/Projects/kitty'); from gateway.model_digest import run_digest; events = run_digest(); print(f'Digest complete: {len(events)} events')</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>55</integer>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Users/jacobbrizinski/Projects/kitty</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/jacobbrizinski/Projects/kitty/logs/model_digest.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jacobbrizinski/Projects/kitty/logs/model_digest.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
