# kitty_gateway/com.kitty.morning-brief.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kitty.morning-brief</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/jacobbrizinski/Projects/kitty/venv/bin/python</string>
        <string>/Users/jacobbrizinski/Projects/kitty/scripts/brief.py</string>
        <string>--notify</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <!-- Load the .env file by sourcing it — brief.py uses python-dotenv -->

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/jacobbrizinski/Projects/kitty/logs/morning-brief.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/jacobbrizinski/Projects/kitty/logs/morning-brief.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
